# misc/debian/etc/cron.d/wansd
*/10 * * * * root  pidof wansd || /usr/local/sbin/wansd
